# networkd-dispatcher.service
[Unit]
Description=Dispatcher daemon for systemd-networkd

[Service]
Type=simple
ExecStart=/usr/bin/networkd-dispatcher
